# templates/superset.service.epp
<%- | 
  String $user,
| -%>
[Unit]
Description="Apache Superset"
After=network.target

[Service]
ExecStart=/bin/superset.gunicorn
User=<%= $user %>
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
