# templates/superset.service.epp
<%- | 
| -%>
[Unit]
Description="Apache Superset"
After=network.target

[Service]
ExecStart=/bin/superset.gunicorn
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
